# videoprocessor.service
[Unit]
Description=Video Processor Service
After=network.target

[Service]
Type=simple
ExecStart=/path/to/deploy/videoprocessor
WorkingDirectory=/path/to/deploy
Restart=always
RestartSec=5s
Environment=CONFIG_PATH=/path/to/deploy/config.json

[Install]
WantedBy=multi-user.target
